# configworker.service
[Unit]
Description=Config Worker Test
 
[Service]
Type=notify
ExecStart=/usr/sbin/configworker/ConfigWorker

[Install]
WantedBy=multi-user.target